# config/launchd/com.saasdashboard.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.saasdashboard.backend</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>{DASHBOARD_PATH}/backend/server.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{DASHBOARD_PATH}/backend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>3000</string>
    <key>HOME</key>
    <string>{HOME_PATH}</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>{DASHBOARD_PATH}/logs/backend.log</string>

  <key>StandardErrorPath</key>
  <string>{DASHBOARD_PATH}/logs/backend.log</string>
</dict>
</plist>
